# .config/systemd/user/nvim-proxy.socket
[Unit]
Description = Socket activation for nvim headless server

[Socket]
ListenStream=6666

[Install]
WantedBy=sockets.target
